# packages/sx05re/emulationstation/system.d/emustation.service
[Unit]
Description=EmulationStation emulator frontend
ConditionPathExists=/var/lock/start.games

[Service]
Environment=HOME=/storage
ExecStartPre=/usr/bin/emustation-config
ExecStartPre=/usr/bin/killkodi.sh
ExecStart=/usr/bin/emulationstation
KillMode=process
TimeoutStopSec=3
Restart=on-failure
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
